# board/recalbox/fsoverlay/etc/udev/rules.d/99-sixpair.rules
#DRIVER=="usb", SUBSYSTEM=="usb", ATTR{idVendor}=="054c", ATTR{idProduct}=="0268", RUN+="/recalbox/scripts/udev-sixpair.sh"

