# board/batocera/fsoverlay/etc/profile.d/20-alias.sh
# ---- ALIAS VALUES ----
alias mc='mc -xc'
alias batocera-check-updates='batocera-upgrade --check-upgrade'
alias which='command -v'
