# board/batocera/fsoverlay/etc/profile.d/20-alias.sh
# ---- ALIAS VALUES ----
alias mc='mc -x'
alias batocera-check-updates='batocera-es-swissknife --update'

