# board/batocera/fsoverlay/etc/profile.d/28-display.sh
export XAUTHORITY=/var/lib/.Xauthority
d="$(getLocalXDisplay)"
[ -n "$d" ] && export DISPLAY="$d"
