# board/batocera/fsoverlay/etc/profile.d/50-system.sh
#!/bin/sh

export MESA_SHADER_CACHE_DIR="/userdata/system/.cache"
