# board/batocera/fsoverlay/etc/profile.d/90-sourcebashrc.sh
if ! [ -f $HOME/.profile ] && ! [ -f $HOME/.bash_profile ]; then
  [ -f $HOME/.bashrc ] && . $HOME/.bashrc
fi
