# package/batocera/utils/wii-u-gc-adapter/99-wii-u-gc-adapter.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="057e", ATTRS{idProduct}=="0337", RUN+="/usr/bin/wii-u-gc-adapter-daemon"
