# board/batocera/allwinner/h700/fsoverlay/etc/udev/rules.d/99-wifi-powersave.rules
# Prevent the RTw88 from going into power save mode
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan0", RUN+="/usr/sbin/iw dev wlan0 set power_save off"

