# board/batocera/rockchip/rk3128/fsoverlay/etc/udev/rules.d/80-headset.rules
SUBSYSTEM=="switch", ACTION=="add|change", ENV{SWITCH_NAME}=="Headset", RUN+="/usr/bin/headset-hotplug.sh"
