# package/batocera/utils/xow/99-xow.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="02e6", RUN+="/usr/bin/xow-daemon"
