# package/batocera/core/batocera-scripts/rules/80-switch-screen.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/bash /usr/bin/batocera-switch-screen-checker"
