# package/batocera/core/batocera-scripts/rules/80-switch-screen.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/bin/bash /usr/bin/batocera-switch-screen-checker-delayed"
