# package/batocera/core/batocera-scripts/rules/99-controller-event.rules
SUBSYSTEM=="input", ENV{ID_INPUT_JOYSTICK}=="1", ACTION=="add|remove", RUN+="/usr/bin/batocera-restart-evmapy"
